# com.stockboard.trading.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent to auto-start StockBoard (backend + frontend + Tailscale Funnel)
  at login. Credentials are read by start-stockboard.sh from .env.trading.

  Install:
    cp com.stockboard.trading.plist ~/Library/LaunchAgents/
    launchctl load  ~/Library/LaunchAgents/com.stockboard.trading.plist
  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.stockboard.trading.plist

  Note: requires TWS paper logged in with API on. Edit the path below if your
  checkout is not at ~/StockBoard.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.stockboard.trading</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd "$HOME/StockBoard" &amp;&amp; ./start-stockboard.sh --funnel</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/stockboard.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/stockboard.err.log</string>
</dict>
</plist>
